# profile.d/aem.sh
#!/usr/bin/env bash

################################################################################
# AWS Environment Manager - https://github.com/bonusbits/aws_env_manager (Will replace aws.sh scripting)
################################################################################
# if [ ! -h "/usr/local/bin/aem" ]; then
# 	ln -s "/Users/levon/Development/github/bonusbits/aws_env_manager/aem.sh" /usr/local/bin/aem
# fi
