# profile.d/aliases.sh
#!/usr/bin/env bash

################################################################################
# BASH Aliases
################################################################################
#alias ls="ls"
alias l="ls"
alias rm="rm -i"
alias mv="mv -i"
alias rsync="rsync -r --rsh=ssh"
alias psg="ps -ef | grep "
alias cls="clear"
alias h="history | grep -i "
