# profile.d/circleci.sh
#!/usr/bin/env bash

################################################################################
# CircleCI CLI Shortcuts
################################################################################
alias cci="circleci "
alias cci-validate="circleci config validate "
alias cci-process="circleci config process .circleci/config.yml"
alias cci-exec="circleci local execute --job "
