# profile.d/java.sh
#!/usr/bin/env bash

################################################################################
## SET JAVA VERSION
################################################################################
# Macos
# /System/Library/Frameworks/JavaVM.framework/Versions/Current/Commands/java_home -v 1.8.0_201 --exec javac -version
## Path
# /Library/Java/JavaVirtualMachines/jdk1.8.0_201.jdk/Contents/Home
# java_version=1.8.0_201
# export IDEA_JDK=$(/usr/libexec/java_home -v $java_version)
# export RUBYMINE_JDK=$(/usr/libexec/java_home -v $java_version)
# export JDK_HOME=$(/usr/libexec/java_home -v $java_version)
# export JAVA_HOME=$(/usr/libexec/java_home -v $java_version)

# Using SDKMan to maintain JDK now 02/20/2019
