# profile.d/jx.sh
#!/usr/bin/env bash

################################################################################
# Jenkins-X Configuration - https://jenkins-x.io/
################################################################################
export PATH="$PATH:$HOME/.jx/bin"
