# profile.d/packer.sh
#!/usr/bin/env bash

################################################################################
# Environment Variables
################################################################################
export PATH=$PATH:$HOME/packer