# profile.d/terraform.sh
#!/usr/bin/env bash

################################################################################
# Terraform
################################################################################
alias tf="terraform $1"
alias tfa="terraform apply"
alias tfd="terraform destroy"
alias tfs="terraform show"

#PATH=$PATH:$HOME/Applications/terraform