# profile.d/upgrades.sh
#!/usr/bin/env bash

################################################################################
# Upgrades
################################################################################
function upgrade_all (){
    upgrade_homebrew
    upgrade_sdkman
}
